# rootfs/etc/services.d/php7-fpm/run
#!/usr/bin/with-contenv sh

exec php-fpm7 -F -y /etc/php7/php-fpm.conf
